# setup/trackmon.service
[Unit]
Description=Trackmon Server

[Service]
ExecStart=/usr/local/bin/trackmon-server -config=/etc/trackmon-server/trackmonserv.conf

[Install]
WantedBy=multi-user.target
